# scripts/launchd/com.user.rift-snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Periodic rift-layout snapshotter. Runs rift-auto-snapshot every 300 seconds.
  Edit ProgramArguments path if you install rift-auto-snapshot elsewhere.

  Install:
    cp scripts/launchd/com.user.rift-snapshot.plist ~/Library/LaunchAgents/
    # adjust the absolute path inside if needed
    launchctl load ~/Library/LaunchAgents/com.user.rift-snapshot.plist

  Disable:
    launchctl unload ~/Library/LaunchAgents/com.user.rift-snapshot.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.rift-snapshot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>$HOME/.bin/rift-auto-snapshot</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/rift-snapshot.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/rift-snapshot.log</string>
</dict>
</plist>
